# etc/profile.d/ops.csh
setenv OPS_ROOT /home/ops/witzenmann-ops
setenv OPS_LOG /var/log/ops
setenv OPS_RUN /run/ops
setenv OPS_BACKUPS /home/ops/z_backup
